# debian/furios-quirks-regenerate-ssh-keys.furios-regenerate-ssh-keys.service
[Unit]
Description=Regenerates SSH keys at first boot
Before=ssh.service
ConditionPathExists=!/etc/ssh/furios-keys-regenerated

[Service]
Type=oneshot
ExecStart=/usr/bin/furios-regenerate-ssh-keys.sh

[Install]
WantedBy=multi-user.target ssh.service
