# keyring/systemd/user/gnome-keyring-daemon.service.d/90-initial-setup.conf
[Unit]
# Do not start keyring daemon if keyring is not created
ConditionPathExists=%E/keyring-setup-done
